# ambient-color.timer
[Unit]
Description=Run ambient color at login

[Timer]
OnActiveSec=10s

[Install]
WantedBy=timers.target
