# ambient-color.service
[Unit]
Description=Start ambient color at boot

[Service]
Type=oneshot
ExecStart=%h/.local/bin/ambient
